# shadowsocks/etc/crontabs/root
#m h  dom mon dow   command
#45 * * * * /usr/bin/tunnel_keepalive >> /var/log/tunnel.log 2>&1
0 0,2,4,12,18 * * * /usr/bin/ss-reload
30 1 * * 3,5 /etc/init.d/shadowsocks reload
15 3 * * * /usr/bin/update
15 4 * * * /usr/bin/update-actions
0 0 * * 1 /sbin/reboot
